# cloudflare-ddns.service
[Unit]
Description=Cloudflare Dynamic DNS Updater
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Environment=CF_ZONE_ID=your_zone_id_here
Environment=DYN_ID=your_dns_record_id_here
Environment=CF_DNS_TOKEN=your_cloudflare_token_here
ExecStart=/usr/bin/python3 /usr/local/bin/cloudflare-ddns.py
StandardOutput=append:/var/log/cloudflare-ddns.log
StandardError=append:/var/log/cloudflare-ddns.log

[Install]
WantedBy=multi-user.target
